# ass04/82-dnetshik-usb.rules
SUBSYSTEM=="input", SUBSYSTEMS=="usb", ACTION=="add", RUN+=/sbin/insmod /home/dnetshik/little-penguin/ass04/hello.ko
